# deploy/telegram-bridge.service
[Unit]
Description=TeleClaude Telegram -> Claude Code bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=claude
Group=claude
WorkingDirectory=/opt/telegram-bridge
Environment=HOME=/home/claude
EnvironmentFile=/opt/telegram-bridge/.env
ExecStart=/opt/telegram-bridge/.venv/bin/python -m src.bot.main
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
